# samba/samba-service.conf
[program:samba]
command=/usr/local/bin/run-samba
autostart=false
autorestart=true
stdout_syslog=true
stderr_syslog=true
priority=8
startretries=0
