# rsyslog/rsyslog-service.conf
[program:syslog]
command=/usr/sbin/rsyslogd -n
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=1
