# http/http-service.conf
[program:http]
command=/usr/local/bin/run-nginx
autostart=false
autorestart=true
stdout_syslog=true
stderr_syslog=true
priority=5
