# rsync/rsync-service.conf
[program:rsync]
command=/usr/bin/rsync --daemon --no-detach --verbose
autostart=false
autorestart=true
stdout_syslog=true
stderr_syslog=true
priority=5
