# services/dhcp-service.conf
[program:dhcp]
command=/usr/sbin/dhcpd -4 -f --no-pid -cf "%(ENV_DHCP_DIR)s/dhcpd.conf" -lf "%(ENV_DHCP_DIR)s/dhcpd.leases" %(ENV_IFACE)s
autostart=false
# dhcpd writes everything to both stderr and syslog, so don't log anything here
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_redirect=true
priority=3
