# services/http-service.conf
[program:http]
command=/usr/sbin/nginx -g "daemon off;"
autostart=false
autorestart=true
stdout_syslog=true
stderr_syslog=true
priority=5
