# services/tftp-service.conf
[program:tftp]
command=/usr/sbin/in.tftpd --listen --foreground --verbosity 4 --secure --ipv4 "%(ENV_TFTP_DIR)s"
autostart=false
stdout_syslog=true
stderr_syslog=true
priority=4
startretries=0
