# tftp/tftp-service.conf
[program:tftp]
command=/usr/local/bin/run-tftpd
autostart=false
stdout_syslog=true
stderr_syslog=true
priority=4
