# tftp/tftp-service.conf
[program:tftp]
command=/usr/sbin/in.tftpd --listen --foreground -vvvv --secure --ipv4 /data/tftp
autostart=false
stdout_syslog=true
stderr_syslog=true
priority=4
