# startup-service/startup-service.conf
[program:startup]
command=/usr/local/bin/startup
autostart=true
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_redirect=true
priority=2
