# startup-service/startup-service.conf
[program:startup]
command=/usr/local/bin/startup
autostart=true
startretries=0
autorestart=false
stdout_syslog=true
stderr_syslog=true
priority=2
